# ansible/roles/k3s_dr_backup/templates/k3s-dr-backup.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Daily k3s DR backup to NFS

[Timer]
OnCalendar={{ k3s_dr_backup_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
