# ansible/roles/maintenance/templates/k3s-image-prune.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Weekly prune of unused container images for k3s

[Timer]
OnCalendar={{ k3s_image_prune_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
